# claude-proxy.service
[Unit]
Description=Claude Code CLI Local Proxy
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.openclaw/skills/claude-proxy
ExecStart=%h/.local/bin/uv run --no-project \
  --with fastapi \
  --with "uvicorn[standard]" \
  --with pydantic \
  --with python-dotenv \
  python proxy.py
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=%h/.openclaw/skills/claude-proxy/.env

[Install]
WantedBy=default.target
